# ops/systemd/automem-consolidate.service
[Unit]
Description=automem daily consolidation
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/automem/app
EnvironmentFile=/opt/automem/app/.env
ExecStart=/opt/automem/.venv/bin/python /opt/automem/app/scripts/scheduled_consolidate.py
