# ops/systemd/automem-consolidate.service
[Unit]
Description=automem daily consolidation
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/automem
EnvironmentFile=/opt/automem/backend/.env
ExecStart=/opt/automem/.venv/bin/python /opt/automem/scripts/scheduled_consolidate.py
